# deploy/systemd/backend.service
[Unit]
Description=Jaden Razo Backend Services
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Project-Website
ExecStart=/bin/bash -c "/root/Project-Website/deploy/start-backend.sh"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jadenrazo-backend
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 